# meta-tinyghettobox/recipes-tinyghettobox/system-setup/files/systemd/tgb.system-setup.service
[Unit]
Description=System setup service

[Service]
ExecStart=/usr/bin/tinyghettobox/system-setup.sh
Type=oneshot
RemainAfterExit=false

[Install]
WantedBy=sysinit.target

