# meta-tinyghettobox/recipes-core/psplash/files/psplash-shutdown.service
[Unit]
Description=Psplash Shutdown Service
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/psplash-shutdown -a 180
RemainAfterExit=true

[Install]
WantedBy=shutdown.target
